# hdsleep-dpkg/package/etc/cron.d/hdsleep
# Run hdsleep every minute
*/10 *	* * *	root 	/usr/bin/hdsleep/hdsleep
* * 	* * *	root	/usr/bin/hdsleep/remlock
