# hdsleep-dpkg/package/etc/cron.d/hdsleep
# Run hdsleep every minute
0,30 * 	* * * 	root 	/usr/bin/hdsleep/hdsleep
* *	* * *	root	/usr/bin/hdsleep/remlock
