# hdsleep-dpkg/package/etc/cron.d/hdsleep
# Run hdsleep every minute
* * * * * root /usr/bin/hdsleep/hdsleep
