# deluge-dpkg/package/etc/cron.d/deluge
* * 	* * *	root	/var/lib/deluged/keepalive